# crontabs/test.cron.sh
#!/bin/sh
php /kunden/homepages/16/d380803416/htdocs/crontabs/scripts/test.cron.php
